# services/docker/chronograf-service/supervisord.conf
[supervisord]
nodaemon = true

[program:chronograf]
priority = 1
command = bash -c "/usr/bin/chronograf start --port=8083 --basepath /chronograf"

