# services/docker/grafana-service/conf/supervisord.conf
[supervisord]
nodaemon = true

[program:grafana]
priority = 1
command = /etc/init.d/grafana-server start

